# configs/disable-turbo.service
[Unit]
Description=disable turbo on cpu

[Service]
ExecStart=/opt/disable-turbo.sh

[Install]
WantedBy=multi-user.target
